# deploy/systemd/migi-user.service
[Unit]
Description=Migi self-hosted agent notification server
Documentation=file:%h/project/migi/docs/deployment.md
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=5min
StartLimitBurst=5

[Service]
Type=simple
UMask=0077

EnvironmentFile=%E/migi/migi.env
StateDirectory=migi
StateDirectoryMode=0700
ConfigurationDirectory=migi
ConfigurationDirectoryMode=0700
WorkingDirectory=%S/migi
LoadCredential=server.crt:%E/migi/server.crt
LoadCredential=server.key:%E/migi/server.key

ExecStart=%h/.local/libexec/migi/migi-server \
    -listen ${MIGI_LISTEN} \
    -ingest-listen ${MIGI_INGEST_LISTEN} \
    -agent-listen=${MIGI_AGENT_LISTEN} \
    -admin-listen ${MIGI_ADMIN_LISTEN} \
    -public-endpoint=${MIGI_PUBLIC_ENDPOINT} \
    -agent-endpoint=${MIGI_AGENT_ENDPOINT} \
    -db %S/migi/migi.db \
    -cert %d/server.crt \
    -key %d/server.key

Restart=on-failure
RestartSec=5s
TimeoutStopSec=15s
KillSignal=SIGTERM
OOMPolicy=stop

StandardOutput=journal
StandardError=journal
SyslogIdentifier=migi
LogRateLimitIntervalSec=30s
LogRateLimitBurst=1000

MemoryHigh=192M
MemoryMax=256M
TasksMax=64
LimitNOFILE=4096

NoNewPrivileges=yes
PrivateTmp=yes
PrivateMounts=yes
ProtectSystem=strict
ProtectHome=read-only
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
KeyringMode=private
ReadOnlyPaths=%E/migi %h/.local/libexec/migi
ReadWritePaths=%S/migi

[Install]
WantedBy=default.target
